# aacalc/cron.daily
#!/bin/sh

/home/ubuntu/oaa/aacalc/gen_sample
