# aacalc/cron.hourly
#!/bin/sh

#find /home/ubuntu/oaa.data/static/results \! -path /home/ubuntu/oaa.data/static/results -prune -mmin +60 -exec rm -r \{\} \;
/home/ubuntu/oaa/aacalc/purge_results
