# web/cron.daily
#!/bin/sh

/home/ubuntu/aacalc/aacalc/gen_sample
