# web/cron.daily
#!/bin/sh

/home/ubuntu/aacalc/web/gen_sample
