# web/cron.daily
#!/bin/bash

DIR=/home/ubuntu/aacalc/spia
DATADIR=/home/ubuntu/.spia

mkdir -p $DATADIR/{nominal,real,corporate}

$DIR/fetch_yield_curve.py -t nominal -d $DATADIR/nominal
$DIR/fetch_yield_curve.py -t real -d $DATADIR/real
$DIR/fetch_yield_curve.py -t corporate -d $DATADIR/corporate

#$DIR/gen_sample
