# web/cron.daily
#!/bin/sh

DIR=/home/ubuntu/aacalc/web
DATADIR=/home/ubuntu/aacalc/opal/data

$DIR/fetch_yield_curve -t corporate -d $DATADIR/public/hqm
$DIR/fetch_yield_curve -t real -d $DATADIR/public/rcmt
$DIR/fetch_yield_curve -t nominal -d $DATADIR/public/cmt

#$DIR/gen_sample
